# collector/deploy/systemd/chaingate-collector-npm.service
[Unit]
Description=ChainGate npm collector — one-shot run
Documentation=https://github.com/chaingate
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=vps_admin
Group=vps_admin
WorkingDirectory=/home/vps_admin/projects/chaingate
EnvironmentFile=/home/vps_admin/projects/chaingate/.env
ExecStart=/home/vps_admin/projects/chaingate/.venv/bin/python -m collector.collector npm

# Hard ceilings — match RUN_TIMEOUT_SECONDS in collector.py with some headroom
TimeoutStartSec=900
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

# stdout/stderr go to journald as structured records
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chaingate-collector-npm

# Minimal sandboxing — the collector only needs network + read access to the
# project dir + write access to Postgres (via TCP, so no FS writes to /var/lib)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/vps_admin/projects/chaingate
PrivateTmp=true
