# collector/deploy/systemd/chaingate-collector-pypi.service
[Unit]
Description=ChainGate PyPI collector — one-shot run
Documentation=https://github.com/chaingate
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=chaingate
Group=chaingate
WorkingDirectory=/home/chaingate/projects/chaingate
EnvironmentFile=/home/chaingate/projects/chaingate/.env
ExecStart=/home/chaingate/projects/chaingate/.venv/bin/python -m collector.collector pypi

TimeoutStartSec=900
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

StandardOutput=journal
StandardError=journal
SyslogIdentifier=chaingate-collector-pypi

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/chaingate/projects/chaingate
PrivateTmp=true
